# .config/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /srv/pocketrefriserver/.config/pocketrefriserver.ini

[program:nginx]
command=nginx

