# .config/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /srv/pf_server/.config/pocketrefriserver.ini

[program:nginx]
command=nginx

